# scripts/systemd/nut-server.service.in
[Unit]
Description=Network UPS Tools - power devices information server
After=local-fs.target network.target nut-driver.service
# We don't Require drivers to be successfully started! This would be
# a change of behavior compared to init SysV, and could prevent from
# accessing successfully started, at least to audit a system.
Wants=nut-driver.service
Before=nut-monitor.service

[Service]
ExecStart=@SBINDIR@/upsd
ExecReload=@SBINDIR@/upsd -c reload
Type=forking

[Install]
WantedBy=multi-user.target
