# scripts/systemd/nut-server.service.in
[Unit]
Description=Network UPS Tools - power devices information server
After=local-fs.target network.target nut-driver.service
Requires=nut-driver.service
Before=nut-monitor.service

[Service]
ExecStart=@SBINDIR@/upsd -n
Type=forking

[Install]
WantedBy=multi-user.target
