# scripts/systemd/nut-driver-enumerator.service.in
[Unit]
Description=Network UPS Tools - enumeration of configure-file devices into systemd unit instances
After=local-fs.target
Before=nut-driver.target

[Install]
WantedBy=multi-user.target
