# scripts/systemd/nut-driver.target.in
# Network UPS Tools (NUT) systemd integration
# Copyright (C) 2011-2025 by NUT contributors
# Distributed under the terms of GPLv2+
# See https://networkupstools.org/
# and https://github.com/networkupstools/nut/

[Unit]
Description=Network UPS Tools - target for power device drivers on this system
After=local-fs.target
# network.target
PartOf=nut.target

# TODO: Dedicated man pages?
Documentation=man:ups.conf(5)
Documentation=@NUT_WEBSITE_BASE@/docs/man/ups.conf.html

Documentation=man:nut.conf(5)
Documentation=@NUT_WEBSITE_BASE@/docs/man/nut.conf.html

[Install]
WantedBy=nut.target
