# scripts/systemd/nut-driver.service.in
[Unit]
Description=Network UPS Tools - power device driver controller
After=local-fs.target network.target
StopWhenUnneeded=yes

[Service]
ExecStart=@DRVPATH@/upsdrvctl start
ExecStop=@DRVPATH@/upsdrvctl stop
ExecStart=/sbin/upsdrvctl start
ExecStop=/sbin/upsdrvctl stop
Type=forking

