# scripts/systemd/nut-driver-enumerator.path.in
# Trigger restart of nut-driver-enumerator.service whenever ups.conf is edited

[Path]
PathModified=@CONFPATH@/ups.conf

[Install]
WantedBy=nut.target
