# scripts/systemd/nut-driver.service.in
[Unit]
Description=Network UPS Tools - power device driver controller
After=local-fs.target network.target
StopWhenUnneeded=yes

[Service]
ExecStart=@SBINDIR@/upsdrvctl -n start
ExecStop=@SBINDIR@/upsdrvctl -n stop
Type=forking

