# scripts/systemd/nut-driver.target
[Unit]
Description=Network UPS Tools - power device driver controller
After=local-fs.target network.target
